# deployments/systemd/traffic-analysis-auto-commit.timer
[Unit]
Description=Run Traffic Analysis Platform automatic git commit every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=traffic-analysis-auto-commit.service

[Install]
WantedBy=timers.target
